# graphite/cache/src/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon = true

[program:carbon-cache]
user = www-data
command = /opt/graphite/bin/carbon-cache.py --debug start
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
